# packaging/systemd/flow.service
# Flow — voice dictation, user systemd unit.
#
# Assumes the `flow` entry point is at ~/.local/bin/flow (the default for
# `pipx install`, `uv tool install`, and `pip install --user`).
#
# Install:
#   install -Dm0644 flow.service ~/.config/systemd/user/flow.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now flow.service
#   journalctl --user -u flow -f

[Unit]
Description=Flow — voice dictation (hold-to-talk)
Documentation=https://github.com/lucyfromnaarm/flow-dictation
After=graphical-session.target ydotoold.service
Wants=ydotoold.service
PartOf=graphical-session.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
Environment=YDOTOOL_SOCKET=%t/.ydotool_socket
ExecStart=%h/.local/bin/flow run
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Light hardening — a full sandbox would block PyTorch/whisper JIT and PortAudio.
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
LockPersonality=yes

[Install]
WantedBy=default.target
