# packaging/systemd/ydotoold.service
# ydotoold — user instance, used by Voicisst to type into Wayland windows.
#
# Run as a *user* unit so the socket lives in XDG_RUNTIME_DIR and is owned
# by you (the system-wide ydotoold puts a root-owned socket in /tmp).
# Requires /dev/uinput access — scripts/setup-linux.sh sets up the udev
# rule and 'input' group membership.
#
# Install:
#   install -Dm0644 ydotoold.service ~/.config/systemd/user/ydotoold.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now ydotoold.service

[Unit]
Description=ydotoold (user instance)
Documentation=man:ydotoold(8)

[Service]
Type=simple
ExecStart=/usr/bin/ydotoold --socket-path=%t/.ydotool_socket --socket-own=%U:%U --socket-perm=0600
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
